# systemd/backup.service
[Unit]
Description=snapback backup

[Service]
Type=simple
ExecStart=/home/kbabioch/snapback/backup.sh

