# systemd/backup.service
[Unit]
Description=snapback backup

[Service]
Type=oneshot
ExecStart=/home/kbabioch/snapback/backup.sh

