# systemd/prune.service
[Unit]
Description=snapback prune

[Service]
Type=oneshot
ExecStart=/home/kbabioch/snapback/prune.sh

