# systemd/prune.timer
[Unit]
Description=snapback prune

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target

