# systemd/snapshot.service
[Unit]
Description=snapback snapshot

[Service]
Type=simple
ExecStart=/home/kbabioch/snapback/snapshot.sh

