# systemd/snapback-prune.service
[Unit]
Description=snapback prune

[Service]
Type=simple
ExecStart=/bin/snapback-prune

