# systemd/snapback-prune.timer
[Unit]
Description=snapback prune

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target

