# systemd/snapback-snapshot@.timer
[Unit]
Description=snapback backup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target

