# systemd/snapback-backup.service
[Unit]
Description=snapback backup

[Service]
Type=simple
ExecStart=/bin/snapback-backup

