# systemd/snapshot.timer
[Unit]
Description=snapback backup

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=multi-user.target

